# LUNA-os/config/includes.chroot/etc/xdg/autostart/luna-dock.desktop
[Desktop Entry]
Type=Application
Name=Luna Dock
Exec=python3 /opt/luna-dock/dock.py
NoDisplay=true
X-GNOME-Autostart-enabled=true
