# LUNA-os/config/includes.chroot/etc/skel/.config/autostart/xfce4-panel.desktop
[Desktop Entry]
Type=Application
Name=XFCE Panel
Exec=xfce4-panel
Hidden=true
